# systemd/storjshare.service.in
[Unit]
Description=Storj Share
After=network.target

[Service]
ExecStart=/usr/bin/storjshare start --datadir /home/sevcsik/.storjshare -p PASSPHRASE_PLACEHOLDER
Restart=on-abort
