# systemd/znc.service
[Unit]
Description=ZNC IRC bouncer
After=network.target

[Service]
ExecStart=/usr/bin/znc -f
Restart=on-abort

